# mcp-sharepoint.service
[Unit]
Description=SharePoint MCP Server
After=network.target

[Service]
Type=simple
User=mcp
WorkingDirectory=/opt/mcp-sharepoint
ExecStart=/usr/bin/node src/server.js mcp
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
EnvironmentFile=/opt/mcp-sharepoint/.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mcp-sharepoint/logs

[Install]
WantedBy=multi-user.target